# deployment/power-monitor-scheduler.service
[Unit]
Description=Power Monitor Nightly Scheduler
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/power-monitor
ExecStart=/usr/bin/python3 /opt/power-monitor/scheduler.py --daemon --data-dir /var/www/html --config /opt/power-monitor/config.json
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
